# ha_bot.service
[Unit]
Description=Telegram HA Control Bot
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/ha_bot.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target
